# root/etc/services.d/beets/run
#!/usr/bin/with-contenv bash
umask 002

exec \
	s6-setuidgid abc beet web
