# services/beets/run
#!/bin/bash
umask 002

/sbin/setuser abc beet web
